# Frenchex.Dev.Vos.Lib/Resources/Provisioning/generic/alpine317/4.2.8/v1/docker-ce/install.sh
#!/usr/bin/env bash

sudo apk update
sudo apk add docker docker-compose
sudo rc-update add docker default
